# assignment-3/vagrant/roles/flask/templates/flask_supervisor.conf.j2
[program:app]
command=python {{ app.install_dir }}/flask_app.py
user=root
autostart=true
autorestart=true
stdout_logfile={{ app.log_dir }}/app.log
stderr_logfile={{ app.log_dir }}/app.err.log
